# ros2_ws/host_camera/nuevo-pi-camera-feed.service
[Unit]
Description=Project NUEVO Pi camera feed to V4L2 loopback
After=systemd-modules-load.service
Wants=systemd-modules-load.service
StartLimitIntervalSec=0

[Service]
Type=simple
EnvironmentFile=/etc/default/nuevo-pi-camera
ExecStart=/usr/local/bin/nuevo-pi-camera-feed
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
